# hc_apache/templates/default/httpd.conf.erb
#!/bin/bash
# chkconfig: 345 90 10
# description: Runs the httpd
# processname: httpd
#
# Dynamically generated by Chef on <%= node['fqdn'] %>
#
# Local modifications will be overwritten by Chef.
#
HTTPD_PID=`ps -ef | grep "httpd" |grep '   1'    |awk '{print$2}' `
if [  $# -ne '1' ];then
echo "No Arguements passed ."
echo "Usage :  service httpd {start|stop|status|restart}"
exit 1
fi

StartHttpd(){
        echo "Starting httpd ... "
        bash -c "/usr/local/apache2/bin/apachectl start"
        sleep 1
}

StopHttpd(){
echo "Stopping httpd  "
bash -c "/usr/local/apache2/bin/apachectl stop" ;
sleep 1
echo "httpd Stopped "
}

StatusHttpd(){
bash -c "/usr/local/apache2/bin/apachectl status"
}

case "$1" in
 start)         StartHttpd ;;
 stop)          StopHttpd  ;;
 restart)   StopHttpd;StartHttpd ;;
 status)        StatusHttpd ;;
 *) echo "Incorrect Arguements."
    echo "Usage :  service httpd  {start|stop|restart|status}"
    exit 1
        ;;
esac

