# hc_jasper/templates/product/wig/jasper.erb
#!/bin/bash
# chkconfig: 345 90 10
# description: Runs the Jasper Application Server
# processname: jasper
#
# Dynamically generated by Chef on <%= node['fqdn'] %>
#
# Local modifications will be overwritten by Chef.
#
#export JAVA_HOME=<%= node['hc_java']['java_home'] %>
JASPER_HOME="<%= node['hc_jasper']['wig']['jasper_path']%>"
JASPER_PID=`ps -ef | grep "$JASPER_HOME" | grep java  |awk '{print$2}' `
if [  $# -ne '1' ];then
echo "No Arguements passed ."
echo "Usage :  service jasper {start|stop|status|restart}"
exit 1
fi

StartJasper(){
JASPER_PID=`ps -ef | grep "$JASPER_HOME" | grep java  |awk '{print$2}' `
if [  -z "$JASPER_PID" ];then
        echo "Starting Jasper ... "
        cd $JASPER_HOME ;
        ./ctlscript.sh start &
        sleep 1
        echo "Jasper is running with PID ` ps -ef | grep "$JASPER_HOME" | grep java  |awk '{print$2}' `"
else
        echo "Jasper already running with PID $JASPER_PID"
fi
}

StopJasper(){

JASPER_PID=`ps -ef | grep "$JASPER_HOME" | grep java  |awk '{print$2}' `
if [ ! -z "$JASPER_PID" ];then
echo "Stopping Jasper  "
kill -9 `echo $JASPER_PID` ;
sleep 1
echo "Jasper Stopped "
else
echo "Jasper already stopped"
fi
}

StatusJasper(){
JASPER_PID=`ps -ef | grep "$JASPER_HOME" | grep java  |awk '{print$2}' `
if [ ! -z "$JASPER_PID" ];then
echo "Jasper is running with PID $JASPER_PID"
else
echo "Jasper is not running"
exit 1
fi

}

case "$1" in
 start)         StartJasper ;;
 stop)          StopJasper  ;;
 restart)   StopJasper;StartJasper ;;
 status)        StatusJasper ;;
 *) echo "Incorrect Arguements."
    echo "Usage :  service jasper  {start|stop|restart|status}"
    exit 1
        ;;
esac

